# domapi/service/domapi.service
[Unit]
Description=Domoticz JSON‑Compat Proxy
After=network.target

[Service]
Type=simple

# → Wenn du feste Umgebungsvariablen brauchst, hier setzen
# Environment="PORT=8089" "TARGET=http://DOMOTICZ-IP:PORT"  

ExecStart=/usr/bin/env PORT=8089 TARGET=http://DOMOTICZ-IP:PORT /usr/bin/node /opt/domapi/j.js

# Neustart bei jeder Exit-Code != 0
Restart=on-failure
RestartSec=5s

# Alle Logs verwerfen
StandardOutput=null
StandardError=null

# (optional) Service als nicht‑root laufen lassen
# User=domapi
# WorkingDirectory=/opt/domapi

[Install]
WantedBy=multi-user.target
